# meta/recipes-core/busybox/files/busybox-syslog.service.in
[Unit]
Description=System Logging Service
Wants=busybox-klogd.service
Requires=syslog.socket

[Service]
EnvironmentFile=-@SYSCONFDIR@/default/busybox-syslog
ExecStart=@base_sbindir@/syslogd -n $OPTIONS
StandardOutput=null

[Install]
Also=busybox-klogd.service
