# meta/recipes-core/busybox/files/busybox-klogd.service.in
[Unit]
Description=Kernel Logging Service

[Service]
ExecStart=@base_sbindir@/klogd -n
StandardOutput=null
